# php/logger/env/docker/supervisor/fluentd.ini
[program:fluentd]
command=/data/app/dumper/bin/dumper-worker-php8-jit Dumper LocalDumpJobWorkerQueue dumper-extractor-%(process_num)02d
numprocs=1
process_name=fluentd
killasgroup=true
autorestart=true
user=root
environment=PATH=%(ENV_PATH)s:/root/google-cloud-sdk/bin
redirect_stderr=true
stdout_logfile=/data/log/supervisor/fluentd.log
